# pautomount.service
[Unit]
Description=Python Automatic Mounter

[Service]
PIDFile=/run/pautomount.pid
ExecStart=/usr/bin/pautomount 
Restart=on-abort
